# services/papra/papra.container
[Unit]
Description=Papra - Document Management System
After=papra.network

[Container]
ContainerName=papra
Image=ghcr.io/papra-hq/papra:latest
Network=papra.network
PublishPort=1221:1221

# Data persistence - stores database and documents
Volume=papra-data.volume:/app/app-data

# Application URL - REQUIRED for proper operation
# CHANGE THIS: Set to your actual Papra URL
Environment=APP_BASE_URL=https://papra.internal.crussell.io

# Optional: Environment variables for configuration
# See https://docs.papra.tech/configuration for available options
# Environment=PAPRA_LOG_LEVEL=info
# Environment=PAPRA_MAX_UPLOAD_SIZE=50MB
Environment=TZ=America/New_York

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target



