# services/karakeep-homedash/karakeep-homedash.container
[Unit]
Description=KaraKeep HomeDash - Compact bookmark dashboard
After=karakeep.container
Requires=karakeep.container

[Container]
ContainerName=karakeep-homedash
Image=ghcr.io/codejawn/karakeep-homedash:latest
Network=karakeep.network
PublishPort=8595:8595

# Mount the KaraKeep database (read-only) - shares the same data volume as Karakeep
# The db.db file is stored in the Karakeep data volume at /data/db.db
Volume=karakeep-app-data.volume:/mnt/karakeep-data:ro

# Mount config directory for persistence (dashboard preferences, column order, etc.)
Volume=karakeep-homedash-config.volume:/app/config

# Override entrypoint to create symlink and start server
# HomeDash expects db.db at /app/db.db, but it's at /mnt/karakeep-data/db.db
Exec=/bin/sh -c "ln -sf /mnt/karakeep-data/db.db /app/db.db && python3 server.py"

# Optional: Configure the KaraKeep URL
# CHANGE THIS: Set to your actual KaraKeep URL (e.g., http://k4:3322 or https://karakeep.yourdomain.com)
Environment=KARAKEEP_URL=https://karakeep.internal.crussell.io

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target



