# services/arr/arr-sonarr.container
[Unit]
Description=Sonarr - TV Show Management
After=arr.network

[Container]
ContainerName=sonarr
Image=lscr.io/linuxserver/sonarr:latest
Network=arr.network
PublishPort=8989:8989

# Config persistence
Volume=arr-sonarr-config.volume:/config

# Media library access - MUST match paths in qBittorrent and Jellyfin
# Your NAS is mounted at /mnt/media on the host
Volume=/mnt/media:/media

# Environment variables for LinuxServer image
Environment=PUID=1000
Environment=PGID=1000
Environment=TZ=America/New_York

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target



