# services/arr/arr-prowlarr.container
[Unit]
Description=Prowlarr - Indexer Manager
After=arr.network

[Container]
ContainerName=prowlarr
Image=lscr.io/linuxserver/prowlarr:latest
Network=arr.network
PublishPort=9696:9696

# Config persistence
Volume=arr-prowlarr-config.volume:/config

# Environment variables for LinuxServer image
Environment=PUID=1000
Environment=PGID=1000
Environment=TZ=America/New_York

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target



