# services/arr/arr-jellyfin.container
[Unit]
Description=Jellyfin Media Server
Wants=network-online.target
After=network-online.target arr.network

[Container]
ContainerName=jellyfin
Image=docker.io/jellyfin/jellyfin:latest
Network=arr.network
PublishPort=8096:8096

# GPU passthrough for hardware transcoding
AddDevice=/dev/dri/renderD128:/dev/dri/renderD128

# This fixes the "/etc/mtab" error on some systems
SecurityLabelDisable=true

# Data persistence - uses existing directories
Volume=%h/.config/jellyfin/config:/config
Volume=%h/.config/jellyfin/cache:/cache

# Media library - your NAS
Volume=/mnt/media:/media

# Auto-update to latest image
AutoUpdate=registry

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target



