# cn/services/karakeep/karakeep-meilisearch.container
[Unit]
Description=Karakeep Meilisearch Search Engine
After=karakeep.network

[Container]
ContainerName=karakeep-meilisearch
Image=docker.io/getmeili/meilisearch:v1.10
Network=karakeep.network

# Meilisearch configuration
# CHANGE THIS: Generate a secure master key with: openssl rand -base64 36
Environment=MEILI_MASTER_KEY=IrLHvlILqJDWOufv/wweLQondy1+rq3JufEOXoJg/2sMjkGS

# Mount the Meilisearch data volume
Volume=karakeep-data.volume:/meili_data

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target


