# cn/services/karakeep/karakeep-chrome.container
[Unit]
Description=Karakeep Chrome Browser (for web scraping)
After=karakeep.network

[Container]
ContainerName=karakeep-chrome
Image=gcr.io/zenika-hub/alpine-chrome:124
Network=karakeep.network

# Chrome configuration for headless operation (simplified for Karakeep)
Exec=--no-sandbox --disable-gpu --disable-dev-shm-usage --remote-debugging-address=0.0.0.0 --remote-debugging-port=9222 --hide-scrollbars

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target


