# services/memos/memos.container
[Unit]
Description=Memos - A privacy-first, lightweight note-taking service
After=memos.network

[Container]
ContainerName=memos
Image=docker.io/neosmemo/memos:stable
Network=memos.network
PublishPort=5230:5230

# Data Persistence - SQLite database and assets
Volume=memos-data.volume:/var/opt/memos

# Optional: Database Configuration
# By default, Memos uses SQLite stored in /var/opt/memos
# To use PostgreSQL or MySQL, uncomment and configure:
# Environment=MEMOS_DRIVER=postgres
# Environment=MEMOS_DSN=postgresql://user:password@host:port/dbname

# Optional: Server Configuration
# Environment=MEMOS_MODE=prod
# Environment=MEMOS_ADDR=0.0.0.0
# Environment=MEMOS_PORT=5230

[Service]
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=multi-user.target default.target



